# shared/rootfs/etc/systemd/system/xinit.service
[Unit]
Description=startx
[Service]
ExecStart=/usr/bin/startx
[Install]
WantedBy=graphical.target
